# experms.service
[Unit]
Description=experms

[Service]
Type=simple
ExecStart=/usr/bin/experms

[Install]
WantedBy=multi-user.target
